# base-dev-env/rc.local
#!/bin/bash
echo "Running rc.local"
chown -R tutysara:tutysara /home/tutysara/
